# s6-overlay/s6-rc.d/02-dbus/run
#!/command/with-contenv sh
rm -rf /var/run/dbus/dbus.pid
echo "Starting dbus"
exec dbus-daemon --system --nofork
