# s6-overlay/s6-rc.d/04-snapclient/run
#!/command/with-contenv sh
echo "Starting snapclient"
exec /usr/local/bin/snapclient $ARG
